# deploy/systemd/openclaw-dashboard.service
[Unit]
Description=OpenClaw Multi-Agent Dashboard API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/openclaw-multi-agents
EnvironmentFile=/etc/default/openclaw-multiagent
ExecStart=/usr/bin/python3 -m uvicorn dashboard_api:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 1
Restart=always
RestartSec=5
StartLimitIntervalSec=0
TimeoutStartSec=60
TimeoutStopSec=30
KillMode=process
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-dashboard

[Install]
WantedBy=multi-user.target
