# deploy/systemd/openclaw-multiagent.service
[Unit]
Description=OpenClaw Multi-Agent Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/openclaw-multi-agents
EnvironmentFile=/etc/default/openclaw-multiagent
ExecStart=/usr/bin/python3 /var/www/openclaw-multi-agents/orchestrator.py \
    --task-timeout-sec 1800 \
    --phase-timeout-sec 7200 \
    --retry-attempts 3 \
    --retry-delay-sec 2
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3
TimeoutStartSec=60
TimeoutStopSec=30
KillMode=process
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-multiagent

[Install]
WantedBy=multi-user.target
